# default/systemd/system/getty@tty1.service.d/override.conf
[Service]
Environment=XDG_SESSION_TYPE=x11
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- admin' --noclear --skip-login - $TERM
