# default/etc/udev/rules.d/power_supply.rules
ACTION=="add|change", \
SUBSYSTEM=="power_supply", \
ATTR{online}=="0", \
RUN+="/usr/bin/light -S 50"
